# roles/uptime-kuma/templates/backup/backup.timer.j2
[Unit]
Description=Run {{ service }} backup periodically

[Timer]
OnCalendar={{ cron }}

[Install]
WantedBy=timers.target